# docker/supervisor/supervisord.d/php-fpm.conf
[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm
autostart=true
autorestart=true
numprocs=1
user=www-data
stdout_logfile=/var/www/storage/logs/supervisor/php-fpm-stdout.log
stderr_logfile=/var/www/storage/logs/supervisor/php-fpm-stderr.log
# redirect_stderr=true
